# systemd/pihole-cloudsync.slice
[Unit]
Description=pihole-cloudsync resource limiter slice
Before=slices.target

[Slice]
CPUQuota=50%
